# orin/wavecam/ops/wavecam.service.d/override.conf
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=4

[Service]
Restart=on-failure
RestartSec=5
