# orin/deploy/dashboard.service
# RETIRED legacy :8080 operator dashboard.
# Keep this file only as historical reference. Do not install it; use
# wavecam.service on http://<orin-ip>:8088 instead.
[Unit]
Description=JetsonTracker operator dashboard
After=network-online.target gps-server.service
Wants=network-online.target

[Service]
Type=simple
User=zack
WorkingDirectory=/data/projects/gimbal
ExecStart=/usr/bin/python3 /data/projects/gimbal/dashboard/dashboard.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
